# oss/motion.service
[Unit]
Description=MotionSense

[Service]
Type=forking
ExecStart=/home/pi/openbehavior/oss/motion.py
ExecStop=/usr/bin/pkill motion.py
