# RFID/rfidscanner.service
# systemd unit for starting RFID scanner on boot
[Unit]
Description=RFIDscanner

[Service]
Type=simple
ExecStart=/home/pi/openbehavior/RFID/rfidscanner.py
ExecStop=/usr/bin/pkill rfidscanner.py

[Install]
WantedBy=getty.target
DefaultInstance=tty1
